# deploy/coding-hermes-scheduler.service
[Unit]
Description=Coding Hermes Fleet Scheduler
After=network.target

[Service]
Type=simple
User=kara
ExecStart=/home/kara/coding-hermes-scheduler/coding-herms-scheduler/bin/schedulerd -db /home/kara/.hermes/coding-hermes/scheduler.db -listen 127.0.0.1:9090 --namespace-mode --max-concurrent 12 --min-interval 1m
Restart=always
RestartSec=5
MemoryMax=32G
StandardOutput=journal
StandardError=journal
Environment=HOME=/home/kara
Environment=PATH=/home/kara/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=API_SERVER_KEY=WZJhtlvVMXfW6ydE2hJgVyJpi6TBYN91fZ22vldTqzE

[Install]
WantedBy=default.target
